# crontabs
00 11 * * * php [CONSOLE] digest:send --verbose --force > /dev/null