# crontabs
00 11 * * * php [CONSOLE] digest:send --env=[ENV] --verbose --force >> [ROOT]/app/logs/digests.log